# files/cron.d/updatemirror

59 1 * * *	mirror	/usr/local/bin/update-mirror.sh

