# files/cron.d/updatemirror

59 1 * * *	mirror	/usr/local/bin/update-mirror.sh

30 3 1 * *	root	/usr/sbin/duperemove -dr /data >/var/log/duperemove.log

